# assets/testdaemon.service
[Unit]
Description=Simple testing daemon

[Service]
Type=simple
ExecStart=/usr/local/bin/testdaemon

[Install]
WantedBy=multi-user.target
